# supervisord.conf
[supervisord]
nodaemon=true

[program:tore]
command=/t.sh


[program:rcxstart]
command=/usr/bin/rcxstart


[program:proxplusmikk]
command=/a.sh
