# supervisord.conf
[supervisord]
nodaemon=true

[program:tore]
command=/t.sh

[program:nginxstartx]
command=/e.sh

[program:proxplusmikk]
command=/a.sh
